# misc/gosmee.service
[Unit]
Description=Gosmee Forward

[Service]
ExecStart=gosmee client https://smee.io/abcdef http://localhost:8080

[Install]
WantedBy=default.target
