# misc/gosmee-server.service
[Unit]
Description=Gosmee Server

[Service]
ExecStart=gosmee server
User=gosmee
EnvironmentFile=/etc/gosmee.env

[Install]
WantedBy=default.target
